# atomic-update/systemd/steamos-update-os-plymouth.service.in
#  SPDX-License-Identifier: LGPL-2.1+
#
#  This file is part of steamos-customizations.
#
#  steamos-customizations is free software; you can redistribute it and/or modify
#  it under the terms of the GNU Lesser General Public License as published
#  by the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=SteamOS Update Operating System Plymouth
Wants=steamos-update-os.service rauc.service
Before=steamos-update-os.service rauc.service
After=plymouth-start.service
Conflicts=plymouth-quit.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@sbindir@/steamos-update-os-plymouth

[Install]
WantedBy=steamos-update-os.target
