# misc/systemd/steamos-vmguest.service.in
#  SPDX-License-Identifier: LGPL-2.1+
#
#  Copyright © 2025 Igalia S.L.
#  Copyright © 2025 Valve Corporation.
#
#  This file is part of steamos-customizations.
#
#  steamos-customizations is free software; you can redistribute it and/or
#  modify it under the terms of the GNU Lesser General Public License as
#  published by the Free Software Foundation; either version 2.1 of the License,
#  or (at your option) any later version.

[Unit]
Description=SteamOS VM guest detection
RefuseManualStart=yes
Before=sddm.service
ConditionPathExists=!@sddm_autologin_conf@
ConditionVirtualization=vm

[Service]
Type=oneshot
ExecStart=@libexecdir@/steamos/steamos-set-autologin-session

[Install]
WantedBy=graphical.target
