# misc/systemd/steamos-mkuser.service.in
#  SPDX-License-Identifier: LGPL-2.1+
#
#  This file is part of steamos-customizations.
#
#  steamos-customizations is free software; you can redistribute it and/or modify
#  it under the terms of the GNU Lesser General Public License as published
#  by the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Create SteamOS User
DefaultDependencies=no
After=local-fs.target
Before=shutdown.target sysinit.target
Conflicts=shutdown.target
ConditionFileIsExecutable=@sbindir@/steamos-mkuser
ConditionPathExists=!/home/steamos

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@sbindir@/steamos-mkuser

[Install]
WantedBy=sysinit.target
