# misc/systemd/steamos-cfs-debugfs-tunings.service.in
#  SPDX-License-Identifier: LGPL-2.1+
#
#  This file is part of steamos-customizations.
#
#  Copyright © 2022 Valve Corporation.
#
#  steamos-customizations is free software; you can redistribute it and/or
#  modify it under the terms of the GNU Lesser General Public License as
#  published by the Free Software Foundation; either version 2.1 of the License,
#  or (at your option) any later version.

[Unit]
Description=set multiple CFS tunings in sched debugfs
After=sysinit.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=@libexecdir@/steamos/steamos-cfs-debugfs-settings

[Install]
WantedBy=multi-user.target
