# swap/systemd/remove-old-swap.service.in
#  SPDX-License-Identifier: LGPL-2.1+
#
#  This file is part of steamos-customizations.
#
#  steamos-customizations is free software; you can redistribute it and/or
#  modify it under the terms of the GNU Lesser General Public License as
#  published by the Free Software Foundation; either version 2.1 of the License,
#  or (at your option) any later version.

[Unit]
Description=One-time: Remove old small swapfile if present (for upgrade to larger swap)
DefaultDependencies=no
Before=swapfile.service shutdown.target home-swapfile.swap
Conflicts=shutdown.target
ConditionPathExists=/home/swapfile
RequiresMountsFor=/home/swapfile

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'if [ $(stat -c%s /home/swapfile) -lt $((8*1024*1024*1024)) ]; then swapoff /home/swapfile 2>/dev/null || true; rm -f /home/swapfile; fi'
ExecStartPost=/bin/rm -f /etc/systemd/system/remove-old-swap.service

[Install]
WantedBy=shutdown.target
