# swap/systemd/steamos-rm-default-swapfile.service.in
#  SPDX-License-Identifier: LGPL-2.1+
#
#  Copyright © 2023-2024 Igalia S.L.
#  Copyright © 2023-2024 Valve Corporation.
#
#  This file is part of steamos-customizations.
#
#  steamos-customizations is free software; you can redistribute it and/or
#  modify it under the terms of the GNU Lesser General Public License as
#  published by the Free Software Foundation; either version 2.1 of the License,
#  or (at your option) any later version.

[Unit]
Description=Remove default swapfile created by SteamOS in previous releases
DefaultDependencies=no
Before=multi-user.target
Conflicts=shutdown.target
RequiresMountsFor=/home/swapfile
ConditionPathExists=/home/swapfile
ConditionFileIsExecutable=@libexecdir@/steamos/steamos-rm-default-swapfile

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@libexecdir@/steamos/steamos-rm-default-swapfile /home/swapfile

[Install]
WantedBy=swap.target
