# chainloader/systemd/steamos-boot.service.in
#  SPDX-License-Identifier: LGPL-2.1+
#
#  This file is part of steamos-customizations.
#
#  steamos-customizations is free software; you can redistribute it and/or modify
#  it under the terms of the GNU Lesser General Public License as published
#  by the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=SteamOS Boot Registration
After=graphical.target
Requisite=graphical.target
RefuseManualStart=yes
ConditionFileIsExecutable=@bindir@/steamos-bootconf
ConditionPathExists=/efi/SteamOS
RequiresMountsFor=/efi/SteamOS
Before=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@bindir@/steamos-bootconf --mode booted --output-to input /efi/@bootconf_relpath@

[Install]
WantedBy=graphical.target
