# chainloader/systemd/steamos-poweroff.service.in
#  SPDX-License-Identifier: LGPL-2.1+
#
#  This file is part of steamos-customizations.
#
#  steamos-customizations is free software; you can redistribute it and/or
#  modify it under the terms of the GNU Lesser General Public License as
#  published by the Free Software Foundation; either version 2.1 of the License,
#  or (at your option) any later version.

[Unit]
Description=SteamOS Power Off
After=umount.target
Before=systemd-poweroff.service
DefaultDependencies=no
RefuseManualStart=yes

[Service]
Type=oneshot
RemainAfterExit=yes
StandardOutput=journal+console
ExecStart=@libexecdir@/steamos/steamos-reboot poweroff

[Install]
WantedBy=poweroff.target
