# settings-importer/systemd/steamos-settings-importer.service.in
#  SPDX-License-Identifier: LGPL-2.1+
#
#  This file is part of steamos-customizations.
#
#  steamos-customizations is free software; you can redistribute it and/or modify
#  it under the terms of the GNU Lesser General Public License as published
#  by the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=SteamOS Settings Importer
Wants=local-fs.target
After=local-fs.target
Before=getty@tty1.service getty@tty2.service getty@tty3.service getty@tty4.service getty@tty5.service getty@tty6.service getty@tty7.service NetworkManager.service sddm.service
ConditionPathIsDirectory=/esp/settings

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@sbindir@/steamos-settings-importer

[Install]
WantedBy=getty@tty1.service getty@tty2.service getty@tty3.service getty@tty4.service getty@tty5.service getty@tty6.service getty@tty7.service NetworkManager.service sddm.service
