# glx/systemd/steamos-glx-backend.service.in
#  SPDX-License-Identifier: LGPL-2.1+
#
#  This file is part of steamos-customizations.
#
#  steamos-customizations is Free Software.
#  You can redistribute it and/or modify it under the terms of the
#  GNU Lesser General Public License as published by
#  the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Choose a GLX backend on boot
DefaultDependencies=no
Conflicts=shutdown.target
After=systemd-tmpfiles-setup.service
Before=sysinit.target shutdown.target
RefuseManualStop=yes

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@libexecdir@/steamos/steamos-glx-backend

[Install]
WantedBy=sysinit.target
