# grub/systemd/steamos-mkvarboot.service.in
#  SPDX-License-Identifier: LGPL-2.1+
#
#  This file is part of steamos-customizations.
#
#  steamos-customizations is free software; you can redistribute it and/or modify
#  it under the terms of the GNU Lesser General Public License as published
#  by the Free Software Foundation; either version 2.1 of the License, or
#  (at your option) any later version.

[Unit]
Description=Create SteamOS /var/boot Directory
DefaultDependencies=no
Before=shutdown.target
Conflicts=shutdown.target
RefuseManualStart=yes
ConditionKernelCommandLine=ro
ConditionPathExists=!/var/boot
RequiresMountsFor=/var

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@bindir@/cp -a /boot /var/boot
ExecStart=@bindir@/rm -fr /var/boot/grub
